# install/systemd/nebula-backend.service
[Unit]
Description=NebulaCP Backend (NestJS)
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=nebula
Group=nebula
WorkingDirectory=/usr/local/nebula/apps/backend
EnvironmentFile=/etc/nebula/.env
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nebula-backend

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/nebula /usr/local/nebula

[Install]
WantedBy=multi-user.target
